# launch-agent/com.example.pullrepos.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.pullrepos</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/$(whoami)/Documents/practice-repo/scripts/pull-on-login.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>0</integer>
  <key>StandardOutPath</key>
  <string>/tmp/com.example.pullrepos.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/com.example.pullrepos.err.log</string>
</dict>
</plist>
